# linux/restic-scaffolding/restic-backup.crontab
# /etc/cron.d/restic-backup: crontab entry for system backup via restic
SHELL=/bin/bash
RUNNER_LOG_DIR=/var/log/restic-backup
RUNNER_NTFY_SERVER=ntfy.cdzombak.net
RUNNER_NTFY_TOPIC=dzops
RUNNER_NTFY_ACCESS_TOKEN=

00 01 * * *  root  runner -job-name "restic backup" -ntfy-priority 4 -ntfy-tags x -retries 1 -- /usr/bin/flock -n /tmp/restic-backup.lock ionice -c2 -n7 nice -n19 /etc/restic-backup/backup.sh

00 00 * * *  root  find "$RUNNER_LOG_DIR" -mtime +30 -name "*.log" -delete  >/dev/null
